# config/modmail.service
[Unit]
Description=Modmail

[Service]
User=root
WorkingDirectory=/usr/local/src/modmail
ExecStart=python3 bot.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
