# config/modmail.service
[Unit]
Description=Modmail

[Service]
User=root
WorkingDirectory=/usr/local/src/modmail
ExecStart=/usr/local/bin/pipenv run bot
Restart=on-failure

[Install]
WantedBy=multi-user.target
